# etc/cron.d/cronHeXing
# minute(0-59) hour(0-23) day(1-31) month(1-12) weekday(0-6,0=sunday)
16 9-17 * * 1-5 /bin/bash /usr/local/bin/deliverIP.sh
